# record-and-playback/mconf_decrypter/scripts/mconf_decrypter.slice
[Unit]
Description=Mconf Decrypter Processing
