# record-and-playback/mconf_decrypter/scripts/mconf-decrypter.target
[Unit]
Description=Mconf decrypter processing
Wants=mconf-recording-decrypter.service
StopWhenUnneeded=true
