# record-and-playback/mconf_decrypter/systemd/mconf-recording-decrypter.target
[Unit]
Description=Mconf recording decrypter processing
Wants=mconf-recording-decrypter.service
